# deploy/systemd/step-ca-p11-kit@.container
[Unit]
Description=Step CA - %i
PartOf=step-ca-p11-kit@%i.target
Requires=p11-kit-server@%i.target
After=p11-kit-server@%i.target

[Service]
Restart=always

[Install]
WantedBy=step-ca-p11-kit@%i.target

[Container]
ContainerName=step-ca-%i
Label=app=step-ca
Label=instance=%i
Image=ghcr.io/wtcross/step-ca-p11-kit:latest
Pull=always

UserNS=keep-id:uid=1001,gid=1001
SecurityLabelType=step_ca_policy.process
SecurityLabelLevel=s0:c100,c200
NoNewPrivileges=true
Notify=true

Network=step-ca-p11-kit.network

# PublishPort intentionally omitted - configure via drop-in

Volume=step-ca-p11-kit-%i-data:/home/step/.step
Volume=%t/p11-kit:/run/p11-kit

# Required secrets for step-ca runtime
Secret=hsm-pin-%i,type=mount,target=/run/secrets/hsm-pin
Secret=admin-password-%i,type=mount,target=/run/secrets/admin-password
Secret=root-cert-%i,type=mount,target=/run/secrets/root.crt
Secret=intermediate-cert-%i,type=mount,target=/run/secrets/intermediate.crt

Environment=STEPPATH=/home/step/.step
Environment=STEP_P11KIT_SOCKET_PATH=/run/p11-kit/%i.sock

EnvironmentFile=%h/.config/step-ca/%i.env
